# ansible/roles/sanoid/templates/check-timemachine-snapshots.timer.j2
[Unit]
Description=Hourly Time Machine snapshot check and Discord notify

[Timer]
OnCalendar=*-*-* *:05:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
