# ansible/roles/sanoid/templates/check-pbs-data-snapshots.timer.j2
[Unit]
Description=Hourly PBS data snapshot check and Discord notify

[Timer]
OnCalendar=*-*-* *:10:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
